# ansible/templates/media-automount@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/local/bin/media-automount %I
