# ansible/templates/media-automount@.service
[Unit]
Description=External storage media mounting service on %i

[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/local/bin/media-automount %I
